# HOME/config/systemd/user/desktop.target
[Unit]
Description=Desktop apps target
Requires=graphical.target
After=graphical.target
#Wants=foobar.service

[Install]
