# HOME/config/systemd/user/xinitrc.target
[Unit]
Description=Xinitrc

[Install]
Alias=default.target
